# crontab
*/10 * * * * ~/saythesign/bitcoin_and_reset.sh
